# packages/thorch-bsp/payload/usr/lib/systemd/system/thorch-inputd.service
[Unit]
Description=Thorch hardware hotkey input daemon
After=systemd-udevd.service

[Service]
Type=simple
EnvironmentFile=-/etc/thorch/inputd.conf
ExecStart=/usr/bin/thorch-inputd
Restart=always
RestartSec=1s

[Install]
WantedBy=multi-user.target
