# packages/thorch-bsp/payload/usr/lib/systemd/system/thorch-powerd.service
[Unit]
Description=Thorch power and fake suspend daemon
After=systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/bin/thorch-powerd
Restart=always
RestartSec=1s

[Install]
WantedBy=multi-user.target
